# rpushd.service
[Unit]
Description=rpushd
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/rpushd
Environment=RPUSHD_LISTEN=127.0.0.1:45831
Environment=RPUSHD_SECRET=replace-with-the-subscription-secret
Environment=RPUSHD_PUBLISH_SECRET=replace-with-the-publish-secret
ExecStart=/opt/rpushd/rpushd
Restart=always
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectClock=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallArchitectures=native
MemoryDenyWriteExecute=true
ReadWritePaths=/opt/rpushd

[Install]
WantedBy=multi-user.target
